# scripts/launchd/com.finances.advisor-curate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.finances.advisor-curate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nyhzdev/devroom/battlefields/finances/scripts/cron-advisor-curate.sh</string>
  </array>

  <!-- Domingo 23:30 hora local (Madrid). Weekday 0 = domingo. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>23</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/nyhzdev/.finances/logs/advisor-curate.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nyhzdev/.finances/logs/advisor-curate.log</string>
</dict>
</plist>
