# scripts/launchd/com.finances.advisor-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.finances.advisor-scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nyhzdev/devroom/battlefields/finances/scripts/cron-advisor-scan.sh</string>
  </array>

  <!-- 09:00 y 18:00 hora local (Madrid). launchd coalesce los disparos
       perdidos cuando el Mac despierta. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/nyhzdev/.finances/logs/advisor-scan.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nyhzdev/.finances/logs/advisor-scan.log</string>
</dict>
</plist>
